# supervisord.conf
[supervisord]
nodaemon=true

[program:php5-fpm]
command=service php5-fpm start

[program:nginx]
command=service nginx start
